# socops.service
[Unit]
Description=SOCops — AI-assisted SOC workbench for Wazuh
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/path/to/socops
EnvironmentFile=/path/to/socops/.env
ExecStart=/usr/bin/python3 /path/to/socops/app.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
